# rarog_configs/rpi2/rarog.conf
# Startup all robot nodes

description "bringup Rarog UAV"

console log
respawn

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

env ROS_HOME=/home/sky/.ros
env ROSCONSOLE_CONFIG_FILE=/home/sky/ros/rosconsole.config

setgid sky
setuid sky

script
  . /home/sky/ros/devel/setup.sh
  exec roslaunch rarog_bringup robot_bringup.launch
end script
