# rarog_configs/rpi2/rarog.conf
description "bringup Rarog UAV"

console log
respawn

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

#env ROS_LOG_DIR=/var/log/ros
#env ROS_IP=XXX
env ROS_HOME=/home/sky/ros
env ROSCONSOLE_CONFIG_FILE=$ROS_HOME/rosconsole.config

setgid sky
setuid sky

script
  cd $ROS_HOME
  . ./devel/setup.sh
  exec roslaunch rarog_bringup robot_bringup.launch
end script
